# services/agetty-tty5/run
#!/bin/sh
exec /bin/chpst -P /sbin/agetty 38400 tty5 linux
