# services/agetty-tty6/run
#!/bin/sh
exec /bin/chpst -P /sbin/agetty --autologin leon 38400 tty6 linux
