# services/agetty-tty2/run
#!/bin/sh
exec /bin/chpst -P /sbin/agetty 38400 tty2 linux
