# services/dhcpcd/run
#!/bin/sh
exec dhcpcd -B ${OPTS:=-M} 2>&1
