# services/metalog/run
#!/bin/sh
exec /usr/sbin/metalog > /dev/null
